# portal/deploy/agentchattr-qwen.service
[Unit]
Description=agentchattr Qwen API agent
After=network.target agentchattr.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/agentchattr
EnvironmentFile=/opt/agentchattr/secrets.env
Environment=PATH=/opt/agentchattr/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/agentchattr/.venv/bin/python -m wrapper_api qwen
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
